# config/comic-server.service
# Replace paths before installation. Logs: journalctl -u comic-server -f

[Unit]
Description=Hardgrizz Comics webcomic server
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/polina_site
ExecStart=/path/to/polina_site/scripts/start.sh
# User=www-data
# Environment=POLINA_SITE=/data/polina
# Environment=IMAGE_RESIZE_ENABLED=true
# Environment=IMAGE_RESIZE_MAX_DIM=1600
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
